# solarswarm_setup/system services/iw_dump.service
[Unit]
Description=Collect iw station dump output
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
Restart=no
ExecStart=/usr/local/bin/iw_dump.bash
User=own_name

[Install]
WantedBy=multi-user.target
